# scripts/60-can-conf.rules
SUBSYSTEM=="net", KERNEL=="can0", RUN+="/usr/sbin/ip link set can0 type can bitrate 500000", RUN+="/usr/sbin/ip link set up can0"
